# debian/DEBIAN/postinst
#!/usr/bin/env bash
set -e

. /usr/share/debconf/confmodule

APP_DIR="/opt/agentic-kali"
VENV_DIR="$APP_DIR/.venv"

python3 -m venv "$VENV_DIR"
"$VENV_DIR/bin/pip" install --upgrade pip
"$VENV_DIR/bin/pip" install -e "$APP_DIR"

ln -sf "$VENV_DIR/bin/agentic-kali" /usr/local/bin/agentic-kali
mkdir -p /etc/agentic-kali
db_get agentic-kali/accept-terms
if [ "$RET" != "true" ]; then
  echo "Terms not accepted. Configure package again to proceed."
  exit 1
fi

db_get agentic-kali/azure-endpoint
AZURE_ENDPOINT="$RET"
db_get agentic-kali/azure-key
AZURE_KEY="$RET"
db_get agentic-kali/azure-deployment
AZURE_DEPLOYMENT="$RET"
db_get agentic-kali/azure-api-version
AZURE_API_VERSION="${RET:-2025-04-01-preview}"

if [ ! -f /etc/agentic-kali/config.json ]; then
  cat > /etc/agentic-kali/config.json <<EOF
{
  "TERMS_ACCEPTED": true,
  "AZURE_OPENAI_ENDPOINT": "$AZURE_ENDPOINT",
  "AZURE_OPENAI_API_KEY": "$AZURE_KEY",
  "AZURE_OPENAI_DEPLOYMENT": "$AZURE_DEPLOYMENT",
  "AZURE_OPENAI_API_VERSION": "$AZURE_API_VERSION"
}
EOF
  chgrp sudo /etc/agentic-kali/config.json || true
  chmod 0660 /etc/agentic-kali/config.json
else
  echo "Keeping existing /etc/agentic-kali/config.json"
  chgrp sudo /etc/agentic-kali/config.json || true
  chmod 0660 /etc/agentic-kali/config.json || true
fi

if [ ! -f /etc/agentic-kali/admin.json ]; then
  python3 - <<'PY'
import hashlib, json
from pathlib import Path
path = Path("/etc/agentic-kali/admin.json")
path.write_text(json.dumps({
    "admin_phrase_sha256": hashlib.sha256(b"enable authorized admin mode").hexdigest()
}, indent=2))
PY
fi
chgrp sudo /etc/agentic-kali/admin.json || true
chmod 0660 /etc/agentic-kali/admin.json || true
touch /etc/agentic-kali/scope.json
chgrp sudo /etc/agentic-kali/scope.json || true
chmod 0660 /etc/agentic-kali/scope.json || true

systemctl daemon-reload || true

echo "Run: sudo agentic-kali init /etc/agentic-kali/scope.json"
echo "UI:  sudo systemctl enable --now agentic-kali-ui"

db_stop
